# systemd/nftgeo.timer
[Unit]
Description=Scheduled refresh of the nftgeo geo firewall rules

[Timer]
OnBootSec=2min
OnCalendar=*-*-* 03,15:00:00
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
