# systemd/nftgeo.service
[Unit]
Description=Rebuild the nftgeo geo firewall from rules.conf, country zones, and AbuseIPDB
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nftgeo-update

[Install]
WantedBy=multi-user.target
